# yocto/meta-splash/recipes-core/files/playvideo.service
[Unit]
Description=Play Splash Video

[Service]
ExecStart=/etc/init.d/playvideo.sh
RemainAfterExit=yes

