# yocto/meta-splash/recipes-core/files/playvideo.service
[Unit]
Description=Play Video on Boot

[Service]
ExecStart=/etc/init.d/playvideo.sh
Restart=always

[Install]
WantedBy=multi-user.target

