# yocto/meta-piracer/recipes-core/source/piracercpp.service
[Unit]
Description=piracer-cpp executer
[Service]
ExecStart=/opt/piracer-cpp/gamepad_control
[Install]
WantedBy=multi-user.target
